# .config/systemd/user/hypr-battery-notifier.service
[Unit]
Description=Disable blur and animations when battery is low

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/scripts/battery.sh
